# net.solarnetwork.node.system.ssh/def/solarssh-cleaner.timer
[Unit]
Description=Close discarded SolarSSH sessions daily.

[Timer]
OnCalendar=daily
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
